# config/sys/fwcloud-api.service
[Unit]
Description=fwcloud-api

[Service]
ExecStart=/usr/bin/npm start
ExecStop=/usr/bin/npm run stop
WorkingDirectory=/opt/fwcloud-api/
StandardOutput=/opt/fwcloud-api/logs/stdout.log
StandardError=/opt/fwcloud-api/logs/stderr.log
SyslogIdentifier=fwcloud-api
User=fwcloud
Group=fwcloud
Type=forking

Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
